# docker/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /opt/app/code/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
